# service/systemd/dbus-org.sailfishos.zxing.service
[Unit]
Description=ZXing service

[Service]
Type=dbus
BusName=org.sailfishos.zxing
ExecStart=/usr/bin/zxing-daemon
